# product_aggregator/crontab
#/bin/bash
* * * * * /usr/local/bin/python /app/manage.py check_offers
# An empty line is required at the end of this file for a valid cron file.
